# services/sunshine.service
[Unit]
Description=Sunshine Game Stream Host
After=headless-xorg.service pulseaudio.service
Wants=headless-xorg.service pulseaudio.service

[Service]
Type=simple
User=steam
Environment=DISPLAY=:0
ExecStart=/usr/bin/sunshine /config/sunshine.conf
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/steam

[Install]
WantedBy=multi-user.target
